# dld-mcp-server/data-service/crontab
0 2 * * 1 cd /app && python scraper.py 90 >> /var/log/scraper.log 2>&1