# docker/crontab
0 2 * * 1 cd /app && python scraper.py 90 2>&1 | logger -t scraper